# copy-dir/airi.conf
# AIRi service
#

start on started tty2
stop on runlevel [!2345]

respawn
exec /opt/AIRi/bin/AIRi --nobrowser 80 2>&1 1>/var/log/AIRi.log
